# scripts/storyteller.service
[Unit]
Description=Bedtime Storyteller Service
Documentation=https://github.com/your-repo/bedtime-storyteller
After=network.target
Wants=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/storyteller
Environment=PATH=/opt/storyteller/venv/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONPATH=/opt/storyteller
ExecStart=/bin/bash -c 'source /opt/storyteller/venv/bin/activate && python -m storyteller.main run --daemon'
ExecReload=/bin/kill -HUP $MAINPID
ExecStop=/bin/kill -TERM $MAINPID
KillMode=mixed
TimeoutStopSec=30
# Simplified security settings for initial deployment
NoNewPrivileges=true

# Hardware access for GPIO and audio
DeviceAllow=/dev/gpiomem rw
DeviceAllow=/dev/snd/* rw
DeviceAllow=/dev/dri/* rw
SupplementaryGroups=audio gpio spi i2c

# Resource limits for Pi Zero 2W
MemoryMax=400M
MemorySwapMax=100M
TasksMax=100
LimitNOFILE=1024

# Restart configuration
Restart=always
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=storyteller

[Install]
WantedBy=multi-user.target